# Library/LaunchAgents/org.hammerspoon.Hammerspoon.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>org.hammerspoon.Hammerspoon</string>
    
    <key>ProgramArguments</key>
    <array>
        <string>/Applications/Hammerspoon.app/Contents/MacOS/Hammerspoon</string>
    </array>
    
    <key>EnvironmentVariables</key>
    <dict>
        <key>LANG</key>
        <string>en_US.UTF-8</string>
        <key>LC_ALL</key>
        <string>en_US.UTF-8</string>
    </dict>
    
    <key>RunAtLoad</key>
    <true/>
    
    <key>KeepAlive</key>
    <false/>
    
    <key>ProcessType</key>
    <string>Interactive</string>
    
    <key>StandardOutPath</key>
    <string>/tmp/hammerspoon.out.log</string>
    
    <key>StandardErrorPath</key>
    <string>/tmp/hammerspoon.err.log</string>
</dict>
</plist>
