# Library/LaunchAgents/org.p0deje.Maccy.plist.tmpl
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>org.p0deje.Maccy</string>
    
    <key>ProgramArguments</key>
    <array>
        <string>/Applications/Maccy.app/Contents/MacOS/Maccy</string>
    </array>
    
    <key>EnvironmentVariables</key>
    <dict>
        <key>LANG</key>
        <string>en_US.UTF-8</string>
        <key>LC_ALL</key>
        <string>en_US.UTF-8</string>
    </dict>
    
    <key>RunAtLoad</key>
    <true/>
    
    <key>KeepAlive</key>
    <false/>
    
    <key>ProcessType</key>
    <string>Interactive</string>
    
    <key>StandardOutPath</key>
    <string>{{ .chezmoi.homeDir }}/.local/state/launchagents/maccy/stdout.log</string>
    
    <key>StandardErrorPath</key>
    <string>{{ .chezmoi.homeDir }}/.local/state/launchagents/maccy/stderr.log</string>
    
    <key>Umask</key>
    <integer>77</integer>
</dict>
</plist>
